# systemd/.config/systemd/user/daily-wallpaper.service
[Unit]
Description=Fetch today's wallpaper and apply theme
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/fetch-wallpaper
Nice=10
