# systemd/.config/systemd/user/daily-wallpaper.service
[Unit]
Description=Fetch today's wallpaper and apply theme

[Service]
Type=oneshot
ExecStart=%h/.local/bin/fetch-wallpaper
Nice=10
